# packaging/deb/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create data and log directories
        mkdir -p /var/lib/aeromux
        mkdir -p /var/log/aeromux

        # Create aeromux system user and group (idempotent)
        if ! getent passwd aeromux > /dev/null 2>&1; then
            adduser --system --group --home /var/lib/aeromux --no-create-home aeromux
        fi

        # Add aeromux to plugdev group for USB device access
        usermod -aG plugdev aeromux

        # Set directory ownership
        chown aeromux:aeromux /var/lib/aeromux
        chown aeromux:aeromux /var/log/aeromux

        # Set configuration file ownership
        chown root:aeromux /etc/aeromux/aeromux.yaml

        # Reload udev rules
        udevadm control --reload-rules
        udevadm trigger

        # Reload systemd daemon
        systemctl daemon-reload

        # Display post-install message on fresh install only
        if [ -z "$2" ]; then
            # Check if sudo is available (present on Raspberry Pi OS, not on pure Debian)
            HAS_SUDO=false
            if command -v sudo > /dev/null 2>&1; then
                HAS_SUDO=true
            fi

            echo "================================================================================"
            echo "Aeromux installed successfully!"
            echo "================================================================================"
            echo ""
            echo "Next steps:"
            echo ""
            echo "  1. Edit the configuration:"

            if [ "$HAS_SUDO" = true ]; then
                echo "     sudo nano /etc/aeromux/aeromux.yaml"
            else
                echo "     su -c 'nano /etc/aeromux/aeromux.yaml'"
            fi

            echo ""
            echo "     - Configure your SDR device(s) under 'sdrSources'"
            echo "     - Configure your Beast TCP input sources under 'beastSources' (if any)"
            echo "     - Set your receiver location under 'receiver'"
            echo ""
            echo "  2. (Optional) Download the aircraft database:"

            if [ "$HAS_SUDO" = true ]; then
                echo "     sudo -u aeromux aeromux database update --config /etc/aeromux/aeromux.yaml"
            else
                echo "     su -s /bin/bash -c 'aeromux database update --config /etc/aeromux/aeromux.yaml' aeromux"
            fi

            echo ""
            echo "  3. Enable and start the service:"

            if [ "$HAS_SUDO" = true ]; then
                echo "     sudo systemctl enable --now aeromux"
            else
                echo "     su -c 'systemctl enable --now aeromux'"
            fi

            echo ""
            echo "  4. View logs:"

            if [ "$HAS_SUDO" = true ]; then
                echo "     sudo journalctl -u aeromux -f"
            else
                echo "     su -c 'journalctl -u aeromux -f'"
            fi

            echo ""
            echo "Documentation: https://github.com/nandortoth/aeromux"
            echo "================================================================================"
        fi
        ;;
esac

exit 0
